# kelly-app-v2/com.kellyapp.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kellyapp.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "/Users/rodrigobermudez/projects/new Kelly App/kelly-app-v2/backend" && source venv/bin/activate && python -m uvicorn main:app --host 0.0.0.0 --port 3026</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rodrigobermudez/projects/new Kelly App/kelly-app-v2/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/rodrigobermudez/projects/new Kelly App/kelly-app-v2/backend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rodrigobermudez/projects/new Kelly App/kelly-app-v2/backend.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
